# root/etc/services.d/node/run
#!/usr/bin/with-contenv bash

#
#	Run node app
#
cd /app
exec s6-setuidgid $CONTAINER_USER $NODE_RUN_CMD
